# yggdrasil_api/ftth_ygg.service
[Unit]
Description=FTTH Yggdrasil API Service
After=network.target yggdrasil.service

[Service]
Type=simple
User=aaa
Group=aaa
WorkingDirectory=/home/aaa/fibra/yggdrasil_api
EnvironmentFile=/home/aaa/fibra/yggdrasil_api/.env
ExecStart=/home/aaa/fibra/venv/bin/python -m uvicorn main:app --host ${YGG_HOST} --port ${YGG_PORT}
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
